# init/systemd/ServiSesentry.service
[Unit]
Description=ServiSesentry monitoring check
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/ServiSesentry/main.py
