# init/systemd/ServiSesentry.service
[Unit]
Description=ServiSesentry monitoring service
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Long-running monitoring service: ``--monitor`` runs the scheduler continuously
# at the interval configured in the panel (monitoring → timer_check). No systemd
# .timer is needed; the service loops on its own. (``--monitor -t 0`` would run a
# single pass and exit, but that mode is not used here.)
ExecStart=/usr/bin/python3 /opt/ServiSesentry/main.py --monitor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
