# init/openrc/init.d/ServiSesentry
#!/sbin/openrc-run
# ServiSesentry monitoring daemon
# Configuration: /etc/conf.d/ServiSesentry

: "${SS_PYTHON:=/usr/bin/python3}"
: "${SS_APP:=/opt/ServiSesentry/main.py}"
: "${SS_ARGS:=-d -c}"

command="${SS_PYTHON}"
command_args="${SS_APP} ${SS_ARGS}"
command_background=true
# shellcheck disable=SC2154  # RC_SVCNAME is set by OpenRC at runtime
pidfile="/var/run/${RC_SVCNAME}.pid"

depend() {
	use net
	need localmount
}
